# supervisor.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod
autorestart=true

[program:pooling]
directory = /usr/src/app
command = node src/index.js
autostart = true
autorestart = false